# rootfs/etc/services.d/zabbix-agent/run
#!/usr/bin/with-contenv bash

/usr/sbin/zabbix_agentd -c /etc/zabbix/zabbix_agentd.conf